# nercone-webserver-autoupdator.timer
[Unit]
Description=Nercone Web Server Updator

[Timer]
OnCalendar=*-*-* 00:00:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
